# gunicron.conf
[program:gunicorn]
command=venv/bin/gunicorn kzntree.wsgi:application --bind 0.0.0.0:8000
directory=/kzntree
user=admin
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
